# jarvis_plant_health_system/crontab
*/5 * * * * python /main.py >> /logfile
