# jarvis_plant_health_system/crontab
*/5 * * * * python3 /home/pi/Desktop/jarvisPlantHealth/jarvis_plant_health_system/src/main.py >> /home/pi/Desktop/jarvisPlantHealth/jarvis_plant_health_system/logfile
